# etc/supervisor/conf.d/launch.conf
[program:noVNC]
command=/usr/share/novnc/utils/novnc_proxy --vnc 0.0.0.0:5901 --listen 6901 --cert /etc/ssl/self.pem
user=root
stdout_logfile=/var/log/novnc.log
redirect_stderr=true
