# etc/supervisor/conf.d/runn.conf
[program:vncd]
command=bash -c '/dockerstartup/go_go.sh'

autorestart=true
directory=/root